# root/etc/s6-overlay/s6-rc.d/10-set-timezone/up
/etc/s6-overlay/scripts/10-set-timezone
